# ansible/roles/fixpi/files/etc/profile.d/showcpu.sh
grep -E "(Hardware|Model)" /proc/cpuinfo
