# ansible/roles/wssh/files/etc/systemd/system/wssh.socket
[Unit]
Description=webssh daemon

[Socket]
ListenStream=/run/wssh.sock

[Install]
WantedBy=sockets.target
